# scripts/debian/postrm
#!/bin/bash

# Remove binary
rm /bin/coraza-waf
# Disable service
rm /etc/init.d/coraza-waf
update-rc.d coraza-waf defaults
# Fix inittab
# Remove user